# templates/nvidia-powertarget.timer.j2
[Unit]
Description=Set NVIDIA power limit on boot

[Timer]
OnBootSec=40

[Install]
WantedBy=timers.target
